# wireguard-tools/wireguard-tools.footprint
drwxr-xr-x	alexis/alexis	etc/
drwx------	alexis/alexis	etc/wireguard/
drwxr-xr-x	alexis/alexis	usr/
drwxr-xr-x	alexis/alexis	usr/bin/
-rwxr-xr-x	alexis/alexis	usr/bin/wg
-rwxr-xr-x	alexis/alexis	usr/bin/wg-quick
drwxr-xr-x	alexis/alexis	usr/lib/
drwxr-xr-x	alexis/alexis	usr/lib/systemd/
drwxr-xr-x	alexis/alexis	usr/lib/systemd/system/
-rw-r--r--	alexis/alexis	usr/lib/systemd/system/wg-quick@.service
-rw-r--r--	alexis/alexis	usr/lib/systemd/system/wg-quick.target
drwxr-xr-x	alexis/alexis	usr/share/
drwxr-xr-x	alexis/alexis	usr/share/bash-completion/
drwxr-xr-x	alexis/alexis	usr/share/bash-completion/completions/
-rw-r--r--	alexis/alexis	usr/share/bash-completion/completions/wg
-rw-r--r--	alexis/alexis	usr/share/bash-completion/completions/wg-quick
drwxr-xr-x	alexis/alexis	usr/share/man/
drwxr-xr-x	alexis/alexis	usr/share/man/man8/
-rw-r--r--	alexis/alexis	usr/share/man/man8/wg.8.gz
-rw-r--r--	alexis/alexis	usr/share/man/man8/wg-quick.8.gz
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/dns-hatchet/
-rwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/dns-hatchet/apply.sh
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/dns-hatchet/hatchet.bash
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/dns-hatchet/README
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/embeddable-wg-library/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/embeddable-wg-library/Makefile
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/embeddable-wg-library/README
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/embeddable-wg-library/test.c
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/embeddable-wg-library/wireguard.c
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/embeddable-wg-library/wireguard.h
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/go/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/go/main.go
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/haskell/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/haskell/package.yaml
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/haskell/Setup.hs
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/haskell/src/
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/haskell/src/Data/
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/haskell/src/Data/Time/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/haskell/src/Data/Time/TAI64.hs
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/haskell/src/Main.hs
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/haskell/stack.yaml
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/python/
-rwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/python/main.py
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/rust/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/rust/Cargo.toml
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/rust/src/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/external-tests/rust/src/main.rs
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/extract-handshakes/
-rwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/extract-handshakes/extract-handshakes.sh
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/extract-handshakes/Makefile
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/extract-handshakes/offset-finder.c
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/extract-handshakes/README
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/extract-keys/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/extract-keys/config.c
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/extract-keys/extract-keys.c
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/extract-keys/Makefile
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/extract-keys/README
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/highlighter/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/highlighter/fuzz.c
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/highlighter/gui/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/highlighter/gui/highlight.cpp
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/highlighter/gui/highlight.pro
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/highlighter/highlight.c
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/highlighter/highlighter.c
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/highlighter/highlighter.h
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/highlighter/Makefile
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/highlighter/README
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/json/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/json/README
-rwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/json/wg-json
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/keygen-html/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/keygen-html/keygen.html
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/keygen-html/README
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/keygen-html/wireguard.js
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/launchd/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/launchd/com.wireguard.wg0.plist
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/launchd/README
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/nat-hole-punching/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/nat-hole-punching/nat-punch-client.c
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/nat-hole-punching/nat-punch-server.c
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/nat-hole-punching/README
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/ncat-client-server/
-rwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/ncat-client-server/client-quick.sh
-rwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/ncat-client-server/client.sh
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/ncat-client-server/README
-rwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/ncat-client-server/server.sh
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/reresolve-dns/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/reresolve-dns/README
-rwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/reresolve-dns/reresolve-dns.sh
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/sticky-sockets/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/sticky-sockets/README
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/sticky-sockets/sticky-sockets.c
drwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/synergy/
-rw-r--r--	alexis/alexis	usr/share/wireguard-tools/examples/synergy/README
-rwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/synergy/synergy-client.sh
-rwxr-xr-x	alexis/alexis	usr/share/wireguard-tools/examples/synergy/synergy-server.sh
