# ostree/ostree.footprint
drwxr-xr-x	alexis/alexis	etc/
drwxr-xr-x	alexis/alexis	etc/grub.d/
lrwxrwxrwx	alexis/alexis	etc/grub.d/15_ostree -> /usr/libexec/libostree/grub2-15_ostree
drwxr-xr-x	alexis/alexis	etc/ostree/
drwxr-xr-x	alexis/alexis	etc/ostree/remotes.d/
drwxr-xr-x	alexis/alexis	usr/
drwxr-xr-x	alexis/alexis	usr/bin/
-rwxr-xr-x	alexis/alexis	usr/bin/ostree
-rwxr-xr-x	alexis/alexis	usr/bin/rofiles-fuse
drwxr-xr-x	alexis/alexis	usr/include/
drwxr-xr-x	alexis/alexis	usr/include/ostree-1/
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-async-progress.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-autocleanups.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-blob-reader.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-bootconfig-parser.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-content-writer.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-core.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-deployment.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-diff.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-dummy-enumtypes.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-gpg-verify-result.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-kernel-args.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-mutable-tree.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-ref.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-remote.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-repo-deprecated.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-repo-file.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-repo-finder-avahi.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-repo-finder-config.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-repo-finder.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-repo-finder-mount.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-repo-finder-override.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-repo.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-repo-os.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-sepolicy.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-sign-ed25519.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-sign.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-sysroot.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-sysroot-upgrader.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-types.h
-rw-r--r--	alexis/alexis	usr/include/ostree-1/ostree-version.h
drwxr-xr-x	alexis/alexis	usr/lib/
drwxr-xr-x	alexis/alexis	usr/libexec/
drwxr-xr-x	alexis/alexis	usr/libexec/libostree/
-rwxr-xr-x	alexis/alexis	usr/libexec/libostree/grub2-15_ostree
drwxr-xr-x	alexis/alexis	usr/lib/girepository-1.0/
-rw-r--r--	alexis/alexis	usr/lib/girepository-1.0/OSTree-1.0.typelib
-rwxr-xr-x	alexis/alexis	usr/lib/libostree-1.so.1.0.0
lrwxrwxrwx	alexis/alexis	usr/lib/libostree-1.so.1 -> libostree-1.so.1.0.0
lrwxrwxrwx	alexis/alexis	usr/lib/libostree-1.so -> libostree-1.so.1.0.0
drwxr-xr-x	alexis/alexis	usr/lib/ostree/
-rwxr-xr-x	alexis/alexis	usr/lib/ostree/ostree-prepare-root
-rwxr-xr-x	alexis/alexis	usr/lib/ostree/ostree-remount
drwxr-xr-x	alexis/alexis	usr/lib/pkgconfig/
-rw-r--r--	alexis/alexis	usr/lib/pkgconfig/ostree-1.pc
drwxr-xr-x	alexis/alexis	usr/lib/systemd/
drwxr-xr-x	alexis/alexis	usr/lib/systemd/system/
drwxr-xr-x	alexis/alexis	usr/lib/systemd/system-generators/
-rwxr-xr-x	alexis/alexis	usr/lib/systemd/system-generators/ostree-system-generator
-rw-r--r--	alexis/alexis	usr/lib/systemd/system/ostree-boot-complete.service
-rw-r--r--	alexis/alexis	usr/lib/systemd/system/ostree-finalize-staged-hold.service
-rw-r--r--	alexis/alexis	usr/lib/systemd/system/ostree-finalize-staged.service
-rw-r--r--	alexis/alexis	usr/lib/systemd/system/ostree-prepare-root.service
-rw-r--r--	alexis/alexis	usr/lib/systemd/system/ostree-remount.service
-rw-r--r--	alexis/alexis	usr/lib/systemd/system/ostree-state-overlay@.service
drwxr-xr-x	alexis/alexis	usr/lib/tmpfiles.d/
-rw-r--r--	alexis/alexis	usr/lib/tmpfiles.d/ostree-tmpfiles.conf
drwxr-xr-x	alexis/alexis	usr/share/
drwxr-xr-x	alexis/alexis	usr/share/bash-completion/
drwxr-xr-x	alexis/alexis	usr/share/bash-completion/completions/
-rw-r--r--	alexis/alexis	usr/share/bash-completion/completions/ostree
drwxr-xr-x	alexis/alexis	usr/share/gir-1.0/
-rw-r--r--	alexis/alexis	usr/share/gir-1.0/OSTree-1.0.gir
drwxr-xr-x	alexis/alexis	usr/share/ostree/
drwxr-xr-x	alexis/alexis	usr/share/ostree/trusted.gpg.d/
-rw-r--r--	alexis/alexis	usr/share/ostree/trusted.gpg.d/README-gpg
